# internal/generator/templates/ha-postgres-health.service.tmpl
[Unit]
Description=PostgreSQL HA Health Check for replica {{ .Replica.Host }}
Documentation=https://github.com/HasithDeAlwis/ha-syncgen
After=postgresql.service network.target
Requires=postgresql.service
PartOf=postgresql.service

[Service]
Type=oneshot
User=postgres
Group=postgres
ExecStart={{ .ReplicaDir }}/health_check.sh
WorkingDirectory={{ .ReplicaDir }}
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ha-syncgen-health-check

# Environment variables for the health check script
Environment=PGUSER={{ .Primary.ReplicationUser }}
Environment=PGDATABASE=postgres
Environment=REPLICA_HOST={{ .Replica.Host }}
Environment=PRIMARY_HOST={{ .Primary.Host }}

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/ha-syncgen {{ .Primary.DataDirectory }}

# Restart behavior
RemainAfterExit=no
TimeoutStartSec=60
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
